# vfio/udev/80-usb-libvirt.rules
ACTION=="add", \
    SUBSYSTEM=="usb", \
    RUN+="/home/rholk/sbin/vmusbautoadd"
